# pkg_scripts/sprokest.conf
description "hello-sprokest"

expect fork

start on runlevel [2345]
stop on runlevel [!2345]

exec java -jar /opt/sprokest/sprokest*.jar 2>&1 > /var/log/sprokest/hello.log &
